# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=flask run
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/flask_app.log
stderr_logfile=/var/log/flask_app_err.log

[program:react]
command=npm run start
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/react_app.log
stderr_logfile=/var/log/react_app_err.log
